# udev/rules.d/99-zdocked.rules
DEVPATH=="/devices/pci0000:00/0000:00:02.0/drm/card0/card0-DP-4", RUN+="/home/amos/scripts/dock.sh", RUN+="/usr/bin/logger -t AMOSBIRD add"
